# config/systemd/dep-reprocess.timer
[Unit]
Description=DEP Reprocessing

[Timer]
OnCalendar=*-*-* 20:10:00 America/Chicago
Persistent=true

[Install]
WantedBy=timers.target
